# .work/usefull_scripts/systemd/autostart.sh.timer
[Unit]
Description=Runs autostart.sh script from ~/.work/usefull_scripts on boot after 20 seconds

[Timer]
OnBootSec=20 s
Unit=autostart.sh.service
Persistent=true

[Install]
WantedBy=timers.target
